# supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx
autostart=true 
autorestart=true 
stdout_logfile=/dev/stdout 
stderr_logfile=/dev/stderr 
